# crontab
@reboot /usr/sbin/kplex_monitor.d/boot_test.sh

