# crontab
@reboot /home/pi/kplex_monitor/boot_test.sh

